# contrib/raspi-alpine/root/etc/apk/protected_paths.d/lbu.list
+etc/init.d/ruuvi-prometheus
+etc/init.d/bluetooth-bcm
+etc/runlevels/boot/bluetooth-bcm
+etc/runlevels/default/ruuvi-prometheus
+usr/local/bin/ruuvi-prometheus
+root/.ssh
